# systemd-example/textfile_exporter.service
[Unit]
Description=Prometheus Textfile Exporter
Documentation=https://github.com/SckyzO/textfile_exporter
After=network.target

[Service]
Type=simple
User=textfile_exporter
Group=textfile_exporter
ExecStart=/usr/local/bin/textfile_exporter \
  --textfile.directory="/var/lib/textfile_exporter" \
  --web.listen-address=":9014"

# Security Hardening
ProtectSystem=full
PrivateTmp=true
NoNewPrivileges=true
ProtectHome=true
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictRealtime=true

Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target